# infrastructure/binarylane/ops/assist-platform-health.service
[Unit]
Description=Assist Platform container health check
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/assist-healthcheck
